# meta-demo-05/recipes-systemd/files/iptables@.service
[Unit]  
Description=iptables 
After=network.target  
   
[Service]  
Type=simple 
ExecStart=sudo /usr/sbin/iptables-restore < /etc/iptables/iptables.rules
